# content/etc/s6-overlay/s6-rc.d/init-rspamd/run
#!/command/with-contenv bash

### Set Defaults
ADMIN_PASS=${ADMIN_PASS:-"admin"}
CLAMAV_HOST=${CLAMAV_HOST:-"clamav"}
CLAMAV_PORT=${CLAMAV_PORT:-"3310"}
DCC_HOST=${DCC_HOST:-"dcc-svc"}
DCC_PORT=${DCC_PORT:-"10045"}
REDIS_HOST=${REDIS_HOST:-"redis"}
REDIS_PORT=${REDIS_PORT:-"6379"}
REDIS_DB=${REDIS_DB:-"7"}
REDIS_TIMEOUT=${REDIS_TIMEOUT:-"1.0"}
LOG_LEVEL=${LOG_LEVEL:-"info"} # error,warning,info,debug
DEBUG_MODULES=${DEBUG_MODULES:-"gpt"}

log() {
  echo "[init-rspamd] $*"
}

format_debug_modules() {
	local raw="$1"
	local formatted=""
	local item=""

	# Allow passing an already formatted list: ["module1", "module2"]
	if [[ "$raw" =~ ^\[.*\]$ ]]; then
		echo "$raw"
		return 0
	fi

	# Convert comma-separated values into rspamd list syntax.
	IFS=',' read -ra items <<< "$raw"
	for item in "${items[@]}"; do
		item="${item#${item%%[![:space:]]*}}"
		item="${item%${item##*[![:space:]]}}"
		if [[ -n "$item" ]]; then
			if [[ -n "$formatted" ]]; then
				formatted+=", "
			fi
			formatted+="\"$item\""
		fi
	done

	if [[ -z "$formatted" ]]; then
		echo "[]"
		return 0
	fi

	echo "[$formatted]"
}

copy_rspamd_configs() {
  log "Copying rspamd configuration files"
  cd /assets/rspamd || exit 1
  tar cvf - . | tar xvf - -C /etc/rspamd/
  log "Configuration files copied"
}

configure_worker_controller() {
  sed -i "s/password = .*/password = \"${ADMIN_PASS}\";/" /etc/rspamd/local.d/worker-controller.inc
}

configure_logging() {
  local formatted_debug_modules
  sed -i "s/level = .*/level = \"${LOG_LEVEL}\";/" /etc/rspamd/local.d/logging.inc

  formatted_debug_modules=$(format_debug_modules "$DEBUG_MODULES")
  sed -i "s|debug_modules = .*|debug_modules = ${formatted_debug_modules};|" /etc/rspamd/local.d/logging.inc
}

configure_redis() {
  cat <<EOF > /etc/rspamd/local.d/redis.conf
servers = "$REDIS_HOST:$REDIS_PORT";
db = "$REDIS_DB";
write_servers = "$REDIS_HOST:$REDIS_PORT";
read_servers = "$REDIS_HOST:$REDIS_PORT";
timeout = $REDIS_TIMEOUT;
EOF
}

configure_antivirus() {
  sed -i "s/servers = .*/servers = \"${CLAMAV_HOST}:${CLAMAV_PORT}\";/" /etc/rspamd/local.d/antivirus.conf
}

configure_dcc() {
  sed -i "s/servers = .*/servers = \"${DCC_HOST}:${DCC_PORT}\";/" /etc/rspamd/local.d/dcc.conf
}

setup_persistent_storage() {
  if [ ! -d /run/rspamd ]; then
    mkdir -p /run/rspamd
  fi

  chown -R rspamd:rspamd /var/lib/rspamd /etc/rspamd /run/rspamd

  mkdir -p /tmp/state
  touch /tmp/state/rspamd-init
}

main() {
  copy_rspamd_configs
  configure_worker_controller
  configure_logging
  configure_redis
  configure_antivirus
  configure_dcc
  setup_persistent_storage

  log "rspamd modifications completed"
}

# if [ "$PLUGIN_DKIM" = "TRUE" ] || [ "$PLUGIN_DKIM" = "true" ];  then
# 	if [ ! -d /var/lib/rspamd/dkim ]; then
# 		mkdir -p /var/lib/rspamd/dkim
# 	fi
# 	if [ ! -e /var/lib/rspamd/dkim/dkim.key ]; then
# 		rspamadm dkim_keygen -b 2048 -s dkim -k /var/lib/rspamd/dkim/dkim.key > /var/lib/rspamd/dkim/dkim.txt
# 	fi
# 	chmod 440 /var/lib/rspamd/dkim/*
# cat <<EOF> /etc/rspamd/local.d/dkim_signing.conf
# selector = "mail";
# path = "/var/lib/rspamd/dkim/\${selector}.key";
# allow_username_mismatch = true;
# EOF
# 	cp -R /etc/rspamd/local.d/dkim_signing.conf /etc/rspamd/local.d/arc.conf
# fi 

#FEATURES_ENABLED=`set -o posix; set | sort | grep PLUGIN | grep TRUE |sed -e 's/PLUGIN_//g' | sed -e 's/=TRUE//g' | awk -vRS="" -vOFS=', ' '$1=$1' | tr A-Z a-z`
#echo '** [rspamd] Enabling the following features: '$FEATURES_ENABLED

main "$@"
